# cistern-castellarius.service
[Unit]
Description=Cistern Castellarius — aqueduct scheduler
After=network.target
After=cistern-web.service
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Type=simple
ExecStart=/home/lobsterdog/go/bin/ct castellarius start

Restart=always
RestartSec=5

TimeoutStopSec=360
KillMode=mixed
KillSignal=SIGTERM

StandardOutput=append:/home/lobsterdog/.cistern/castellarius.log
StandardError=append:/home/lobsterdog/.cistern/castellarius.log

Environment=HOME=/home/lobsterdog
Environment=PATH=/home/lobsterdog/go/bin:/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
